# stage1/defaults/containerd/containerd
#!/sbin/openrc-run
#
#
command="/usr/local/bin/containerd"
description="Container Runtime Daemon"
pidfile="/var/run/containerd/containerd.pid"

depend() {
	need net
	after firewall
	provide ntp-client ntp-server
	use dns
}

start() {
	local daemon_args="--background"
	if yesno "$FAST_STARTUP"; then
		ARGS="${ARGS}"
	fi

	ebegin "Starting containerd"
	start-stop-daemon $daemon_args --start --quiet \
		--exec "${command}" \
		--pidfile "${pidfile}" \
		-- ${ARGS}
	eend $? "Failed to start containerd"
}

stop() {
	ebegin "Stopping containerd"
	start-stop-daemon --stop --quiet \
		--pidfile "${pidfile}"
	eend $? "Failed to stop containerd"
}
