# orbit
#!/sbin/openrc-run
#
#
command="/usr/local/bin/orbit"
description="Orbit Agent"
pidfile="/var/run/orbit/orbit.pid"

depend() {
	need net
	after containerd
	provide ntp-client ntp-server
	use dns
}

start() {
	local daemon_args="--background"
	if yesno "$FAST_STARTUP"; then
		ARGS="${ARGS}"
	fi

	ebegin "Starting orbit"
	start-stop-daemon $daemon_args --start --quiet \
		--exec "${command}" \
		--pidfile "${pidfile}" \
		-- ${ARGS}
	eend $? "Failed to start orbit"
}

stop() {
	ebegin "Stopping orbit"
	start-stop-daemon --stop --quiet \
		--pidfile "${pidfile}"
	eend $? "Failed to stop orbit"
}
