# orbit
#!/sbin/openrc-run
#
#
command="/usr/local/bin/orbit-server"
description="Orbit Agent"
pidfile="/var/run/orbit/orbit.pid"
args="--logger /usr/local/bin/orbit-syslog"

depend() {
	need net
	after containerd
	use dns
}

start() {
	local daemon_args="--background"
	ARGS=$args
	if yesno "$FAST_STARTUP"; then
		ARGS="${ARGS}"
	fi

	ebegin "Starting orbit"
	start-stop-daemon $daemon_args --start --quiet \
		--exec "${command}" \
		--stdout-logger "logger -t orbit -p daemon.info" \
		--stderr-logger "logger -t orbit -p daemon.err" \
		--pidfile "${pidfile}" \
		-- ${ARGS}
	eend $? "Failed to start orbit"
}

stop() {
	ebegin "Stopping orbit"
	start-stop-daemon --stop --quiet \
		--pidfile "${pidfile}"
	eend $? "Failed to stop orbit"
}
