# apps/containerd/containerd
#!/sbin/openrc-run
#
supervisor=supervise-daemon

command="/usr/local/bin/containerd"
command_args=""
description="Container Runtime Daemon"
pidfile="/var/run/containerd.pid"

depend() {
	need net
	after firewall
	use dns
}
