# os/buildkit.service
[Unit]
Description=buildkit
After=containerd.service network.target

[Service]
ExecStart=/usr/local/bin/buildkitd --containerd-worker=true --oci-worker=false --addr tcp://0.0.0.0:9500
Restart=always

[Install]
WantedBy=multi-user.target
