# 18.10/node_exporter.service
[Unit]
Description=Prometheus node exporter

[Service]
ExecStart=/usr/local/bin/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
