# stage1/defaults/node_exporter/nodeexporter
#!/sbin/openrc-run
#
#
command="/usr/local/bin/node_exporter"
description="Prometheus Node Metrics"
pidfile="/var/run/node_exporter/node_exporter.pid"
args='--web.listen-address=":9900"'

depend() {
	need net
	use dns
}

start() {
	local daemon_args="--background"
	ARGS=$args
	if yesno "$FAST_STARTUP"; then
		ARGS="${ARGS}"
	fi

	ebegin "Starting node_exporter"
	start-stop-daemon $daemon_args --start --quiet \
		--exec "${command}" \
		--pidfile "${pidfile}" \
		-- ${ARGS}
	eend $? "Failed to start node_exporter"
}

stop() {
	ebegin "Stopping node_exporter"
	start-stop-daemon --stop --quiet \
		--pidfile "${pidfile}"
	eend $? "Failed to stop node_exporter"
}
